# pyfll/packages/packages.d/base-init.postinst
#!/bin/sh -e

say() {
	echo "I: $@..."
}

do_postinst() {
    # OBACHT - muss vor dem Release raus
    systemctl enable debug-shell
}

case "${1}" in
	postinst)
		do_postinst
		;;
	*)
		echo "Usage: ${0} postinst"
		;;
esac

:
