# pyfll/packages/packages.d/base-init.postinst
#!/bin/sh -e

say() {
	echo "I: $@..."
}

do_postinst() {
    # say 'Disable systemd-cron in live-mode'
    # systemctl disable cron.target
    echo
}

case "${1}" in
	postinst)
		do_postinst
		;;
	*)
		echo "Usage: ${0} postinst"
		;;
esac

:
