# pyfll/packages/packages.d/base-init.postinst
#!/bin/sh -e

say() {
	echo "I: $@..."
}

do_postinst() {
    # OBACHT - muss vor dem Release raus
    echo "debug-shell not enabled"
    # systemctl enable debug-shell

    echo "systemd-coredump disabled"
    mkdir -p /etc/systemctl.d
    ln -sf /dev/null /etc/sysctl.d/50-coredump.conf

}

case "${1}" in
	postinst)
		do_postinst
		;;
	*)
		echo "Usage: ${0} postinst"
		;;
esac

:
